# lowpower/idle-dim.service.template
[Unit]
Description=Gambit: dim DSI backlight when kiosk is idle
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Environment=IDLE_TIMEOUT_SECONDS=%IDLE_TIMEOUT_SECONDS%
Environment=DIM_SCRIPT=%DIM_SCRIPT%
ExecStart=/usr/local/bin/gambit-input-idle
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
